# deploy/supervisor_backend.conf
[program:qsl-backend-1]
directory=/opt/qsl-tracker/backend
command=/opt/qsl-tracker/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
user=qsl
autostart=true
autorestart=true
stdout_logfile=/opt/qsl-tracker/logs/backend_1.log
stderr_logfile=/opt/qsl-tracker/logs/backend_1.error.log
environment=PATH="/opt/qsl-tracker/venv/bin",ENVIRONMENT="production"

[program:qsl-backend-2]
directory=/opt/qsl-tracker/backend
command=/opt/qsl-tracker/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002
user=qsl
autostart=true
autorestart=true
stdout_logfile=/opt/qsl-tracker/logs/backend_2.log
stderr_logfile=/opt/qsl-tracker/logs/backend_2.error.log
environment=PATH="/opt/qsl-tracker/venv/bin",ENVIRONMENT="production"

[program:qsl-backend-3]
directory=/opt/qsl-tracker/backend
command=/opt/qsl-tracker/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8003
user=qsl
autostart=true
autorestart=true
stdout_logfile=/opt/qsl-tracker/logs/backend_3.log
stderr_logfile=/opt/qsl-tracker/logs/backend_3.error.log
environment=PATH="/opt/qsl-tracker/venv/bin",ENVIRONMENT="production"
